# scripts/systemd/aether-history.service
[Unit]
Description=Aether Edge history service (aether-history)
PartOf=aether.target

[Service]
Type=simple
EnvironmentFile=/etc/aether/aether.env
Environment=API_HOST=127.0.0.1
ExecStart=/opt/aether/bin/aether-history
Restart=on-failure
RestartSec=2
User=root
WorkingDirectory=/opt/aether

[Install]
WantedBy=aether.target
